# sd/zkreplay-status.timer
[Unit]
Description=Make public status of zkreplay-process

[Timer]
# Every 2m
OnCalendar=*-*-* *:0/2:00

[Install]
WantedBy=multi-user.target
