# sd/zkreplay-status.timer
[Unit]
Description=Make public status of zkreplay-process

[Timer]
# Every 10m
OnCalendar=*-*-* *:0/10:00

[Install]
WantedBy=multi-user.target
